# src/flicker/flicker.service
[Unit]
Description=Flicker

[Service]
ExecStart=/srv/bob/flicker.py
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
